# cron/crontab
0 0 */2 * * /var/www/html/bin/console app:ingest-ongoing
